# wxingest.service
# /etc/systemd/system/wxingest.service
[Unit]
Description=Weather Window NBM Ingestion
After=network.target

[Service]
Type=oneshot
User=gpetty
WorkingDirectory=/home/gpetty/WxApp
Environment=DATA_DIR=/12TB2/NBM
Environment=AWS_NO_SIGN_REQUEST=yes
ExecStart=/home/gpetty/WxApp/.venv/bin/python -m backend.app.ingest --postprocess
ExecStartPost=-/usr/bin/curl -s -X POST http://127.0.0.1:8001/admin/reload
StandardOutput=journal
StandardError=journal
